# src/offline/mactahoe-liquid-kde-theme.service
[Unit]
Description=MacTahoe Liquid KDE — auto light/dark theme switcher
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/mactahoe-theme-switch watch
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
